# backend/shieldmendai-backend.service
[Unit]
Description=ShieldMendAI read-only backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/ShieldMendAi
EnvironmentFile=-/root/ShieldMendAi/backend/.env
Environment=SHIELDMEND_HOST=127.0.0.1
Environment=PORT=8787
ExecStart=/usr/bin/python3 /root/ShieldMendAi/backend/server.py
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
